# config/systemd/books-hardcover-sync.service
[Unit]
Description=Books Hardcover request sync
Wants=network-online.target
After=network-online.target books-calibre.service

[Service]
Type=oneshot
EnvironmentFile=/etc/books/books.env
ExecStart=/opt/books/bin/books-hardcover-sync sync
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ReadWritePaths=/srv/books
