# config/systemd/books-hardcover-sync.service
[Unit]
Description=Books Hardcover request sync
Wants=network-online.target
After=network-online.target books-calibre.service

[Service]
Type=oneshot
EnvironmentFile=/etc/books/books.env
WorkingDirectory=/home/exedev/books
Environment=NODE_OPTIONS=--no-warnings
ExecStart=/home/exedev/books/scripts/books hardcover sync
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ReadWritePaths=/srv/books
