# config/systemd/books-hardcover-sync.timer
[Unit]
Description=Run Books Hardcover request sync every five minutes

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min
AccuracySec=30s
Persistent=true

[Install]
WantedBy=timers.target
