# config/systemd/books-portal.service
[Unit]
Description=Books portal
After=network.target

[Service]
User=books
Group=books
EnvironmentFile=/etc/books/books.env
ExecStart=/opt/books/bin/books-portal
Restart=on-failure
RestartSec=3
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
