# config/systemd/books-calibre.service
[Unit]
Description=Books Calibre content server
Documentation=https://manual.calibre-ebook.com/server.html
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=books
Group=books
EnvironmentFile=/etc/books/books.env
ExecStart=/usr/bin/calibre-server \
    --listen-on 127.0.0.1 \
    --port ${CALIBRE_PORT} \
    --enable-auth \
    --auth-mode=basic \
    --userdb ${BOOKS_USERDB} \
    --log ${BOOKS_LOG_DIR}/calibre.log \
    --access-log ${BOOKS_LOG_DIR}/calibre-access.log \
    --max-log-size 10 \
    --book-list-mode details_list \
    --max-opds-items 100 \
    --max-opds-ungrouped-items 100 \
    ${BOOKS_LIBRARY_DIR}
Restart=on-failure
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=full
ReadWritePaths=/srv/books

[Install]
WantedBy=multi-user.target
