# config/systemd/books-node.service
[Unit]
Description=Books Node app
After=network.target

[Service]
Environment=NODE_OPTIONS=--no-warnings
EnvironmentFile=/etc/books/books.env
WorkingDirectory=/home/exedev/books
ExecStart=/usr/bin/node /home/exedev/books/src/server.js
Restart=on-failure
RestartSec=3
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ReadWritePaths=/srv/books

[Install]
WantedBy=multi-user.target
